# cfg/left4dead2server.service
[Unit]
Description=left4dead2server
After=network.target

[Service]
Type=simple
User=dev
ExecStart=/home/dev/docker-left4dead2server/cfg/run.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target

